# deploy/homepage-status.service
[Unit]
Description=Homepage server status collector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/homepage-status
Environment=SERVER_STATUS_CONFIG=/etc/homepage-status/config.json
ExecStart=/usr/bin/python3 /opt/homepage-status/server_status_service.py
Restart=always
RestartSec=3
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
